# scripts/python/launchd/com.txotracker.fetchmarket.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  本機 launchd 排程 — 抓大盤行情寫 Firestore (Fubon SDK 版)。

  setup.sh --launchd 會把 __REPO_ROOT__ / __VENV_PYTHON__ 替換成實際路徑後
  複製到 ~/Library/LaunchAgents/ 並 load。

  排程:
    11:00 / 13:25 / 14:05 (週一至週五,TW)
    + RunAtLoad=true → 開機立即跑一次 (補抓昨天/今天)

  electrical: launchd 自己會在錯過排程時間 (筆電休眠/關機) 嘗試在恢復後執行,
              但 macOS 14+ 對非互動排程的 catch-up 不保證 — RunAtLoad 是更穩的補抓機制。
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.txotracker.fetchmarket</string>

  <key>ProgramArguments</key>
  <array>
    <string>__VENV_PYTHON__</string>
    <string>__REPO_ROOT__/scripts/python/fetch_market_fubon.py</string>
    <string>--slot</string>
    <string>auto</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key><integer>1</integer>
      <key>Hour</key><integer>11</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>2</integer>
      <key>Hour</key><integer>11</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>3</integer>
      <key>Hour</key><integer>11</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>4</integer>
      <key>Hour</key><integer>11</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>5</integer>
      <key>Hour</key><integer>11</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>1</integer>
      <key>Hour</key><integer>13</integer>
      <key>Minute</key><integer>25</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>2</integer>
      <key>Hour</key><integer>13</integer>
      <key>Minute</key><integer>25</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>3</integer>
      <key>Hour</key><integer>13</integer>
      <key>Minute</key><integer>25</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>4</integer>
      <key>Hour</key><integer>13</integer>
      <key>Minute</key><integer>25</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>5</integer>
      <key>Hour</key><integer>13</integer>
      <key>Minute</key><integer>25</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>1</integer>
      <key>Hour</key><integer>14</integer>
      <key>Minute</key><integer>5</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>2</integer>
      <key>Hour</key><integer>14</integer>
      <key>Minute</key><integer>5</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>3</integer>
      <key>Hour</key><integer>14</integer>
      <key>Minute</key><integer>5</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>4</integer>
      <key>Hour</key><integer>14</integer>
      <key>Minute</key><integer>5</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>5</integer>
      <key>Hour</key><integer>14</integer>
      <key>Minute</key><integer>5</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>__REPO_ROOT__/scripts/python/.launchd.log</string>

  <key>StandardErrorPath</key>
  <string>__REPO_ROOT__/scripts/python/.launchd.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
